# etc/supervisor/conf.d/aphlict.conf
[program:aphlict]
command=sudo -u www /srv/phabricator/bin/aphlict start --foreground
