# etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -c /srv/etc/nginx/nginx.conf
