# etc/supervisor/conf.d/phd.conf
[program:phd]
command=sudo -u www /srv/phabricator/bin/phd start --foreground
