# etc/supervisor/conf.d/sshd_ctrl.conf
[program:sshd-ctrl]
command=/usr/sbin/sshd -D -f /srv/etc/sshd_ctrl/sshd_config
